# docker/django/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/Project/docker/django/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
